# scripts/restart-transmission-container.service
[Unit]
Description="Restart a docker container named Transmission"

[Service]
ExecStart=/bin/sh -c "docker restart transmission"

